# ai.evo.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.evo.bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rameshkrishnan/.nvm/versions/node/v23.5.0/bin/node</string>
        <string>/Users/rameshkrishnan/evo-bridge/server.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rameshkrishnan/evo-bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/rameshkrishnan/evo-bridge/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rameshkrishnan/evo-bridge/bridge.log</string>
</dict>
</plist>
